# vagrant/provisioning/ideato.common/crontab/cocoon
*/15	*	*	*	*	/usr/bin/curl -H "ORA-JWT:eyJhbGciOiJSUzI1NiIsInR5cCI6IkpXUyJ9.eyJ1aWQiOiIwMDAwMDAwMC0wMDAwLTAwMDAtMDAwMC0wMDAwMDAwMDAwMDAiLCJleHAiOiIiLCJpYXQiOiIxNDYxNjgwNjAzIn0.kJ0CzJTgQ5D0dxRQzDF8TnqykDmyH6Xo6IO7TRXIPVQ4bs3TY80IUUYFAzzQRxp6ZHfx4rx8LrnK-Bm24Rbb6vM-8oQQGudVlCVNqR4f8j9yV0fB2cboW5xedzjMP8XTuuJyu93lzo_GoHA3uWQJQ_VrpDs2g5qcovF7Qj2QgfwFwcZQDaAZiLyppstRm4KnjkGJFcXyt5aU-Ad8pIxCDF1wAj1LzUa2YbyJhLisJZAA4zXCyjNH7K8t_rs_DSLHuagNfSB_CFDoo60ckVWPWdS_I9DUBLTU7K2ryXH_9sXyulbF-9ADE2Ax7KFv1rfAJs6WLeDhXvRSDQ8kgtw" --data "" http://cocoon.ideato.it/task-management/tasks/votingresults/idea-items >> /tmp/cocoon-cron-idea-items.log  2>&1
*/15	*	*	*	*	/usr/bin/curl -H "ORA-JWT:eyJhbGciOiJSUzI1NiIsInR5cCI6IkpXUyJ9.eyJ1aWQiOiIwMDAwMDAwMC0wMDAwLTAwMDAtMDAwMC0wMDAwMDAwMDAwMDAiLCJleHAiOiIiLCJpYXQiOiIxNDYxNjgwNjAzIn0.kJ0CzJTgQ5D0dxRQzDF8TnqykDmyH6Xo6IO7TRXIPVQ4bs3TY80IUUYFAzzQRxp6ZHfx4rx8LrnK-Bm24Rbb6vM-8oQQGudVlCVNqR4f8j9yV0fB2cboW5xedzjMP8XTuuJyu93lzo_GoHA3uWQJQ_VrpDs2g5qcovF7Qj2QgfwFwcZQDaAZiLyppstRm4KnjkGJFcXyt5aU-Ad8pIxCDF1wAj1LzUa2YbyJhLisJZAA4zXCyjNH7K8t_rs_DSLHuagNfSB_CFDoo60ckVWPWdS_I9DUBLTU7K2ryXH_9sXyulbF-9ADE2Ax7KFv1rfAJs6WLeDhXvRSDQ8kgtw" --data "" http://cocoon.ideato.it/task-management/tasks/votingresults/completed-items >> /tmp/cocoon-cron-completed-items.log  2>&1