# vagrant/provisioning/ideato.common/crontab/cocoon
5   *  *  *  * centos   /usr/bin/php /var/www/vhosts/welo-backend/public/index.php closepolls idea-items >> /tmp/welo-backend-cron-idea-items.log  2>&1
10  *  *  *  * centos   /usr/bin/php /var/www/vhosts/welo-backend/public/index.php closepolls completed-items >> /tmp/welo-backend-cron-completed-items.log  2>&1
15  0  *  *  * centos   /usr/bin/php /var/www/vhosts/welo-backend/public/index.php reminder >> /tmp/welo-backend-cron-remider.log 2>&1
*/30  *  *  *  * centos  /usr/bin/php /var/www/vhosts/welo-backend/public/index.php sync >> /tmp/welo-backend-cron-kanbanize-sync.log 2>&1
20  0  *  *  * centos  /usr/bin/php /var/www/vhosts/welo-backend/public/index.php closeshares >> /tmp/welo-backend-cron-shares-timebox-sync.log 2>&1
